# drive-stack/carputer/ota-agent-drive-stack.service
[Unit]
Description=Activate drive-stack overlay for the carputer
DefaultDependencies=no
After=local-fs.target
Wants=drive-stack.target
Before=drive-stack.target

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/ota-agent/drive-stack-activate.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
